# freshrss/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: FreshRSS
# Configures NGINX for use with freshrss
# ==============================================================================

declare ingress_entry

# Get ingress entry and strip leading /
ingress_entry="$(bashio::addon.ingress_entry)"
[[ "${ingress_entry}" =~ ^/(.*)$ ]] && ingress_entry="${BASH_REMATCH[1]}"

# Generate Ingress configuration.
bashio::var.json \
    ingress_host "$(bashio::addon.ip_address)" \
    ingress_port "^$(bashio::addon.ingress_port)" \
    ingress_entry "${ingress_entry}" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

# Generate direct access configuration, if enabled.
if bashio::var.has_value "$(bashio::app.port 80)"; then
    if bashio::config.true 'ssl'; then
        bashio::config.require.ssl
    fi
    bashio::var.json \
        certfile "$(bashio::config 'certfile')" \
        keyfile "$(bashio::config 'keyfile')" \
        ssl "^$(bashio::config 'ssl')" \
        | tempio \
            -template /etc/nginx/templates/direct.gtpl \
            -out /etc/nginx/servers/direct.conf
fi
