# rsshub/rootfs/etc/s6-overlay/s6-rc.d/init-rsshub/up
#!/usr/bin/with-contenv bashio

# Post-initialization hook (optional)
bashio::log.debug "RSSHub post-init complete"
exit 0
