# rsshub/rootfs/etc/s6-overlay/s6-rc.d/rsshub/run
#!/usr/bin/with-contenv bash
set -e

# Apply user configuration from addon options
if bashio::has_value "request_retry"; then
    export REQUEST_RETRY="$(bashio::config 'request_retry')"
fi

if bashio::has_value "request_timeout"; then
    export REQUEST_TIMEOUT="$(bashio::config 'request_timeout')"
fi

if bashio::has_value "cache_expire"; then
    export CACHE_EXPIRE="$(bashio::config 'cache_expire')"
fi

if bashio::has_value "cache_content_expire"; then
    export CACHE_CONTENT_EXPIRE="$(bashio::config 'cache_content_expire')"
fi

if bashio::has_value "logger_level"; then
    export LOGGER_LEVEL="$(bashio::config 'logger_level')"
fi

# Set default port and hostname
export PORT=1200
export HOSTNAME=0.0.0.0

bashio::log.info "Starting RSSHub with configured environment variables"

# Start RSSHub via upstream entrypoint
exec /app/docker-entrypoint.sh
