# ddns-go/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

bashio::log.info "Waiting for port 9876 to open..."

# Wait for transmission to become available
bashio::net.wait_for 9876 localhost 9876

bashio::log.info "Starting NGinx..."

exec nginx
