# rsshub/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh

if [ "$1" -ne 0 ] && [ "$1" -ne 256 ]; then
  s6-svscanctl -t /var/run/s6/services
fi
