# rsshub/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start nginx service
# ==============================================================================
exec nginx -g 'daemon off;'
